# copy/etc/systemd/system/bootstrap.service
[Unit]
Description=bootstrap lx-brand image
After=network.target

[Service]
ExecStart=/usr/local/bin/bootstrap

[Install]
WantedBy=default.target
